# files/etc/uci-defaults/99_fans
#!/bin/sh
uci set fancontrol.settings.fan_file='/sys/devices/virtual/thermal/cooling_device1/cur_state'
uci commit fancontrol
exit 0
